# deploy_webapp/gunicorn_process.conf
[supervisord]
[program:gunicorn]
user=runuser
directory=/code/helloworld/
command=gunicorn helloworld:application -w 3 -b 0.0.0.0:8000